# templates/jupyter.service.erb
[Unit]
Description=Jupyter Notebook

[Service]
Type=simple
PIDFile=/run/jupyter-<%= @service_name %>.pid
# Step 1 and Step 2 details are here..
# ------------------------------------
ExecStart=/usr/local/bin/jupyter-notebook --no-browser
User=<%= @group %>
Group=<%= @user %>
WorkingDirectory=/home/<%= @user %>/notebooks
Restart=always
RestartSec=10
#KillMode=mixed

[Install]
WantedBy=multi-user.target
